# quickshell/ii/scripts/systemd/ii-super-overview.service
[Unit]
Description=Super key tap -> Quickshell ii overview toggle
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 %h/.local/bin/ii_super_overview_daemon.py
Restart=on-failure
RestartSec=1

[Install]
WantedBy=default.target
